# main.crontab
* * * * * /usr/bin/python Baykov/app.py

